# conf/supervisor_config_worker.ini
[program:app-celery-worker]
command=celery worker -A app -l INFO
directory=/Users/zhangxin/PycharmProjects/celeryapp
environment=PATH="/Users/zhangxin/PycharmProjects/celeryapp/venv/bin"
stdout_logfile=/Users/zhangxin/PycharmProjects/celeryapp/log/work.log
stderr_logfile=/Users/zhangxin/PycharmProjects/celeryapp/log/work_err.log
autostart=true
priority=999