# conf/supervisor_config_flower.ini
[program:app-celery-flower]
command=celery flower --broker=redis://localhost:6379/1 --address=127.0.0.1 --port=9999 --basic_auth=admin:123456
directory=/Users/zhangxin/PycharmProjects/celeryapp
environment=PATH="/Users/zhangxin/PycharmProjects/celeryapp/venv/bin"
stdout_logfile=/Users/zhangxin/PycharmProjects/celeryapp/log/flower.log
stderr_logfile=/Users/zhangxin/PycharmProjects/celeryapp/log/flower_err.log
autostart=true
priority=777